# execution-cockpit/deploy/cockpit-worker.conf
; Supervisor program for the Execution Cockpit queue worker.
; Install:
;   sudo cp deploy/cockpit-worker.conf /etc/supervisor/conf.d/cockpit-worker.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start cockpit-worker:*
; Check:  sudo supervisorctl status
; Adjust the two paths and `user` to match your server.

[program:cockpit-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/execution-cockpit/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --backoff=10
directory=/var/www/execution-cockpit
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=deploy
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/execution-cockpit/storage/logs/worker.log
stopwaitsecs=3600
